# srcs/influxdb/supervisor.conf
[supervisord]
nodaemon=true

[program:telegraf]
command = sh -c "telegraf && kill 1"

[program:influxdb]
command= sh -c "/usr/sbin/influxd -config /etc/influxdb.conf && kill 1"